# launchagents/com.klody.web-mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.klody.web-mcp</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/klodynlov/Projets/klody-code-ai/scripts/start-web-mcp.sh</string>
        <string>--http</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/klodynlov/Projets/klody-code-ai</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/Users/klodynlov/Library/Logs/klody-web-mcp.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/klodynlov/Library/Logs/klody-web-mcp.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
